# docker/etc/s6-overlay/s6-rc.d/mysql/run
#!/command/execlineb -P
with-contenv

# Give MySQL 5 seconds to start up
s6-notifyoncheck -d -s 5000 \
    /usr/sbin/mysqld
